# scripts/finance-sync.service
[Unit]
Description=Daily Finance Pi sync
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=pi
WorkingDirectory=/home/pi/finance-pi
ExecStart=/usr/bin/docker compose -f compose/docker-compose.yml run --rm app python -m src.orchestrator
