# scripts/finance-sync.timer
[Unit]
Description=Run Finance Pi sync daily at 06:00

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
